# build_app/ai.bam.screenrecord.updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.bam.screenrecord.updater</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Library/Application Support/ScreenRecorder/update_helper.sh</string>
    <string>--scheduled</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>WatchPaths</key>
  <array>
    <string>/Users/Shared/ScreenRecorder_update_now</string>
  </array>
  <key>StandardOutPath</key>
  <string>/Library/Logs/ScreenRecorder/updater.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Library/Logs/ScreenRecorder/updater.launchd.err.log</string>
</dict>
</plist>
